# root/etc/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

@reboot mrminer nohup /root/mrminer/cron/cron.sh > /dev/null 2>&1 &
*/15 * * * *  mrminer bash -c 'sudo /root/mrminer/cron/update.sh' >/dev/null 2>&1

## test 6
