# root/etc/rc.local
#!/bin/sh

nohup /root/mrminer/cron/cron.sh &

bash /root/mrminer/bootstrap.sh &

exit 0
