# root/etc/rc.local
#!/bin/sh
sleep 3

bash /root/mrminer/bootstrap.sh &

exit 0
