# root/etc/rc.local
#!/bin/sh

bash /root/mrminer/bootstrap.sh &

exit 0
